# Code/shm-rpi-based-sensors/interface/rsc/rc.local
#!/bin/sh -e
#
# rc.local /etc/rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0"
#

nohup python3 /home/pi/Desktop/interface/rsc/io.py &

exit 0
